# etc/skel/.config/autostart/plconf.desktop
[Desktop Entry]
Exec=/usr/local/bin/plconf
Icon=dialog-scripts
Name=plconf
Path=
Type=Application
X-KDE-AutostartScript=true
